# services/vestor.service
[Unit]
Description=Vestor LED matrix display
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/vestor
EnvironmentFile=/home/pi/vestor/.env
Environment=PYTHONUNBUFFERED=1
# Entry point is vestor-tracker.py. setcap on the venv python grants
# cap_sys_nice without root; drop_privileges=True is set in display.
ExecStart=/home/pi/vestor/env/bin/python3 vestor-tracker.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
